# SoftwareC++/rc.local
#Place in /etc/rc.local to autostart the software at startup.
sudo /home/pi/softwareC++/motionControler &